# php/supervisord.conf
[supervisord]
nodaemon=true

[program:fpmstatus]
command=start-status-server.sh

[program:fpm]
command=php-fpm

